# deploy/sloane-samehadaku-discover.service
[Unit]
Description=sloane samehadaku daily new-series discovery
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=dwizzy
WorkingDirectory=%h/dwizzyOS/sloane
Environment=PYTHONPATH=%h/dwizzyOS:%h/dwizzyOS-HQ
Environment=ROUTER_BASE_URL=http://192.168.100.6:20128/v1
EnvironmentFile=%h/.config/sloane/ingest.env
ExecStart=%h/dwizzyOS/.venv-adk/bin/python -m sloane.ingest samehadaku --discover
