# deploy/sloane-db-tunnel.service
[Unit]
Description=sloane SSH tunnel to DOS-pgb on homeserver (local 6432)
After=network-online.target
Wants=network-online.target

[Service]
# Autossh-style persistent tunnel: forward local 6432 -> homeserver's
# pgbouncer (127.0.0.1:6432, homeserver-local). Restart on drop so the
# ingest/discover services that Require this always have a live DB path.
ExecStart=/usr/bin/ssh -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -L 6432:127.0.0.1:6432 dwizzy@192.168.100.6
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
