# deploy/sloane-samehadaku-ingest.service
[Unit]
Description=sloane samehadaku 2h feed-delta ingest
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dwizzyOS/sloane
Environment=PYTHONPATH=%h/dwizzyOS:%h/dwizzyOS-HQ
# DOS_PGB_URL (DB password) + ROUTER_API_KEY live in ingest.env (not committed):
EnvironmentFile=%h/.config/sloane/ingest.env
# 9router (LLM for merge fuzzy + enrich) lives on the homeserver too:
Environment=ROUTER_BASE_URL=http://192.168.100.6:20128/v1
ExecStart=%h/dwizzyOS/.venv-adk/bin/python -m sloane.ingest samehadaku
